# launchagents/com.atlas.jarvis-compact-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.atlas.jarvis-compact-weekly</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Template: run scripts/generate_launchagents.sh to substitute {{JARVIS_ROOT}} + {{LOG_DIR}}. -->
        <string>{{JARVIS_ROOT}}/.venv/bin/python</string>
        <string>{{JARVIS_ROOT}}/scripts/run_compaction.py</string>
        <string>--tier</string>
        <string>weekly</string>
    </array>

    <!-- Sunday at 4:00 AM local time (runs AFTER Sunday's daily at 3:00 AM) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{LOG_DIR}}/jarvis-compact-weekly.out</string>
    <key>StandardErrorPath</key>
    <string>{{LOG_DIR}}/jarvis-compact-weekly.err</string>

    <key>WorkingDirectory</key>
    <string>{{JARVIS_ROOT}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
